# form-api/verus-wiki-api.service
[Unit]
Description=Verus Wiki Form Submission API
After=network.target

[Service]
Type=simple
User=cluster
WorkingDirectory=/home/cluster/verus-wiki-retype/form-api
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

# Set your admin token via .env or systemd override: systemctl edit verus-wiki-api
# Generate one with: openssl rand -base64 32
# Environment=ADMIN_TOKEN=<your-token-here>
Environment=PORT=3737

[Install]
WantedBy=multi-user.target
