# form-api/verus-wiki-api.service
[Unit]
Description=Verus Wiki Form Submission API
After=network.target

[Service]
Type=simple
User=cluster
WorkingDirectory=/home/cluster/verus-wiki-retype/form-api
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

# Set your admin token here (generate with: openssl rand -hex 32)
Environment=ADMIN_TOKEN=750dc544fa7ce79c77cb79bd24d3ad9ef6b4677f9e586f3df1f44c9fd83a19db
Environment=PORT=3737

[Install]
WantedBy=multi-user.target
